# roles/stockcheckdeploy/files/stockcheckbuild/build/supervisor-conf/stockcheck_queueelasticsearchcategory.conf
[program:stockcheck_queueelasticsearchcategory]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/html/yii queueelasticsearchcategory/listen --verbose=1 --color=0
autostart=true
autorestart=true
stopasgroup=true
stopsignal=QUIT
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/console/runtime/logs/queueelasticsearchcategory.log
