# roles/stockcheckdeploy/files/stockcheckbuild/build/supervisord.conf
[supervisord]
nodaemon=true
user=root

[include]
files = /etc/supervisor/conf.d/*.conf

[program:init]
directory=/var/www/html
command=/init  ; S6 init will start all services
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/backend/runtime/init.log
